# etc/cron.d/homesensor
*/5 * * * * root    /home/homesensor/homesensor/detector/owfs_temp.pl
